# crontab
0 9 * * * bundle exec runner /rails/lib/batch/remind_event.rb  Batch::RemindEvent.remind_event
