# crontab
0 0 * * * bundle exec rails runner Batch::RemindEvent.remind_event
